# Tests/Service/fixtures/supervisord_config_cli.conf


[program:markup_job_queue_testenv_testqueuea]
command=/usr/local/bin/rabbitmq-cli-consumer -e "/vagrant/app/console markup:job_queue:consumer --strict-exit-code --env=dev --no-debug" -c /etc/rabbitmq-cli-consumer/config/testenv_testqueuea_consumer.conf -V -i --strict-exit-code
stderr_logfile=/vagrant/app/logs/supervisord.error.log
stdout_logfile=/vagrant/app/logs/supervisord.out.log
autostart=false
autorestart=true
stopsignal=QUIT
startsecs=0


[program:markup_job_queue_testenv_testqueueb]
command=/usr/local/bin/rabbitmq-cli-consumer -e "/vagrant/app/console markup:job_queue:consumer --strict-exit-code --env=dev --no-debug" -c /etc/rabbitmq-cli-consumer/config/testenv_testqueueb_consumer.conf -V -i --strict-exit-code
stderr_logfile=/vagrant/app/logs/supervisord.error.log
stdout_logfile=/vagrant/app/logs/supervisord.out.log
autostart=false
autorestart=true
stopsignal=QUIT
startsecs=0


[group:markup_testenv]
programs=markup_job_queue_testenv_testqueuea,markup_job_queue_testenv_testqueueb
